# daemon/mavlink-daemon.service
[Unit]
Description=MAVLink Monitor Daemon
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/mavlink-daemon/mavlink_daemon.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
